# deploy/template/flower.ini
[program:celery-flower]
command = path/to/flower -A celery_entry.celery_app --pidfile=var/pids/flower.pid --loglevel=info
directory = path/to/myce
priority = 1000  ; 让flower最后启动
user = ubuntu
redirect_stderr = true
stdout_logfile = var/logs/celery_flower.log
stdout_logfile_maxbytes = 20MB
stdout_logfile_backups = 5
