# deploy/template/worker.ini
[program:worker]
command = path/to/celery -A celery_entry.celery_app worker -n page_info_worker -Q test-queue --loglevel=info -P gevent -c 1
directory = path/to/myce
priority = 999
user = ubuntu
redirect_stderr = true
stdout_logfile = var/logs/worker.log
stdout_logfile_maxbytes = 20MB
stdout_logfile_backups = 5
